# besu-network/roles/besu-validator3/files/pantheon.service
[Unit]
Description=Ethereum Pantheon client

[Service]
Type=simple
Restart=always
ExecStart=/home/trung/besu-mitx/start-pantheon.sh

[Install]
WantedBy=default.target
