# besu-network-local/roles/besu-bootnode/files/orion.service
[Unit]
Description=Orion client

[Service]
Type=simple
Restart=always
ExecStart=/opt/besu-bootnode/start-orion.sh

[Install]
WantedBy=default.target
