# besu-network-local/roles/boot-node-1/files/orion.service
[Unit]
Description=Orion client

[Service]
Type=simple
Restart=always
ExecStart=/node1/start-orion.sh

[Install]
WantedBy=default.target
