# besu-network-local/roles/validator-node-1/files/orion.service
[Unit]
Description=Orion client

[Service]
Type=simple
Restart=always
ExecStart=/node2/start-orion.sh

[Install]
WantedBy=default.target
